# systemd/netbox-sync-interfaces.timer
[Unit]
Description=IT Nexus 定時同步: LibreNMS Interfaces -> NetBox (每日 02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
