# systemd/netbox-sync-glpi.service
[Unit]
Description=IT Nexus 同步: NetBox -> GLPI
After=network.target postgresql.service netbox-sync-librenms.service
OnFailure=sync-failure-notify@%n.service

[Service]
Type=oneshot
User=netbox
Group=netbox
WorkingDirectory=/opt/netbox/scripts
EnvironmentFile=/opt/netbox/scripts/.env
ExecStart=/opt/netbox/scripts/venv/bin/python3 /opt/netbox/scripts/sync_netbox_to_glpi.py
StandardOutput=journal
StandardError=journal

# 安全硬化 (Security Hardening)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
